# deployments/packages/debian/postinst
#!/bin/bash
set -e

echo "Configuring GPU Health Agent..."

# Setup directories and permissions
mkdir -p /var/lib/gpuhealth
chmod 755 /var/lib/gpuhealth /usr/bin/gpuhealth

echo "Setting up GPU Health Agent service..."

case "$1" in
    configure)
        if [ ! -d /run/systemd/system ]; then
            echo "Systemd not available - manual service management required"
        else
            systemctl daemon-reload || echo "[WARNING] Failed to reload systemd daemon"
            
            if [ -z "$2" ]; then
                # Fresh installation
                echo "Fresh installation - enabling and starting service..."
                systemctl enable gpuhealthd.service || echo "[WARNING] Failed to enable service"
                systemctl start gpuhealthd.service || echo "[WARNING] Failed to start service"
            else
                # Upgrade
                echo "Upgrading from version $2 - restarting service..."
                if systemctl is-active --quiet gpuhealthd.service; then
                    systemctl restart gpuhealthd.service || echo "[WARNING] Failed to restart service"
                else
                    systemctl start gpuhealthd.service || echo "[WARNING] Failed to start service"
                fi
            fi
        fi
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        # Recovery scenarios
        if [ -d /run/systemd/system ] && systemctl is-enabled --quiet gpuhealthd.service; then
            systemctl start gpuhealthd.service || true
        fi
        ;;
esac

# Final status and guidance
echo ""
if [ -d /run/systemd/system ]; then
    if systemctl is-active --quiet gpuhealthd.service; then
        echo "[SUCCESS] GPU Health Agent installed and started successfully!"
        echo ""
        echo "Configuration: sudo vi /etc/default/gpuhealth"
        echo "Logs: sudo journalctl -u gpuhealthd -f"
    else
        echo "[WARNING] GPU Health Agent installed but service is not running"
        echo "Start manually: sudo systemctl start gpuhealthd"
    fi
else
    echo "[SUCCESS] GPU Health Agent installed (systemd not available)"
    echo "Manual start: /usr/bin/gpuhealth run --log-level=info"
fi
echo "" 