# deployments/packages/scripts/rpm/preremove
#!/bin/sh
set -e

service=fleetintd.service

case "$1" in
    0)
        echo "Erasing Fleet Intelligence Agent..."
        if [ -d /run/systemd/system ]; then
            systemctl --no-reload disable --now "$service" >/dev/null || true
        fi
        echo "Fleet Intelligence Agent service stopped and disabled"
        ;;
    1)
        echo "Preparing Fleet Intelligence Agent for upgrade..."
        # Leave the service running during upgrade. postinstall try-restarts it only if it is already active.
        ;;
esac

exit 0
