# deployments/packages/systemd/fleetintd.service
[Unit]
Description=Fleet Intelligence Daemon
Wants=nvidia-dcgm.service
After=network.target local-fs.target nvidia-dcgm.service

[Service]
EnvironmentFile=-/etc/default/fleetint

ExecStart=/usr/bin/fleetint run $FLEETINT_FLAGS
ExecReload=/bin/kill -HUP $MAINPID

Restart=always
RestartSec=5s
TimeoutStartSec=300
TimeoutStopSec=30

User=root
Group=root
LimitNOFILE=40000

[Install]
WantedBy=multi-user.target
