# deployments/packages/scripts/prerm
#!/bin/bash
set -e

case "$1" in
    remove|0)
        echo "Removing Fleet Intelligence Agent..."
        
        if [ -d /run/systemd/system ]; then
            # Stop if active
            if systemctl is-active --quiet fleetintd.service 2>/dev/null; then
                echo "Stopping fleetintd service..."
                systemctl stop fleetintd.service || echo "[WARNING] Failed to stop service"
            fi
            
            # Disable if enabled
            if systemctl is-enabled --quiet fleetintd.service 2>/dev/null; then
                echo "Disabling fleetintd service..."
                systemctl disable fleetintd.service || echo "[WARNING] Failed to disable service"
            fi
        fi
        
        echo "Fleet Intelligence Agent service stopped and disabled"
        ;;
        
    upgrade|1|2)
        echo "Preparing Fleet Intelligence Agent for upgrade..."
        # Don't stop service during upgrade
        ;;
        
    failed-upgrade|deconfigure)
        echo "Fleet Intelligence Agent operation failed, stopping service..."
        if [ -d /run/systemd/system ] && systemctl is-active --quiet fleetintd.service 2>/dev/null; then
            systemctl stop fleetintd.service || true
        fi
        ;;
esac
